# lab7/roles/bind_exporter/templates/bind_exporter.service.j2
[Unit]
Description=Bind Exporter

[Service]
User=root
ExecStart=/opt/bind_exporter/bind_exporter

[Install]
WantedBy=multi-user.target
