# lab7/roles/bind_exporter/templates/bind_exporter.service.j2
[Unit]
Description=Bind Exporter
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/opt/bind_exporter/bind_exporter
Restart=always

[Install]
WantedBy=multi-user.target
