# exam/roles/prometheus/templates/prometheus.service.j2
[Unit]
Description=Prometheus


[Service]
User=prometheus
ExecStart=/usr/bin/prometheus --config.file=/etc/prometheus/prometheus.yml

[Install]
WantedBy=multi-user.target


